# setup/evolver.service
[Unit]
Description=EvoMap Evolver - Continuous Evolution & Asset Distillation
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/evomap-farmer
EnvironmentFile=/root/evomap-farmer/.env
ExecStart=/usr/bin/env evolver run --loop
Restart=always
RestartSec=30
StandardOutput=append:/root/evomap-farmer/logs/evolver-service.log
StandardError=append:/root/evomap-farmer/logs/evolver-service.log

[Install]
WantedBy=multi-user.target